# curator4-911799909.conf
[program:curator4-911799909]
command=python3 /home/curator4/curator4-911799909/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-911799909/test.err.log
stdout_logfile=/home/curator4/curator4-911799909/test.out.log